# packages/deb/git-auto-pull/usr/share/git-auto-pull/git-auto-pull.service
[Unit]
Description=GitLab Webhook auto pull service
After=network.target remote-fs.target

[Service]
Type=forking
PIDFile=/var/run/git-auto-pull.pid
ExecStart=/usr/sbin/git-auto-pull
ExecStop=/bin/kill -WINCH ${MAINPID}
PrivateTmp=true

[Install]
WantedBy=multi-user.target

